# rc.local
#!/bin/bash

DIR="$( cd "$( dirname "${BASH_SOURCE[0]}" )" && pwd )"

mount -t tmpfs -o size=30m tmpfs $DIR/video
mount -t tmpfs -o size=10m tmpfs $DIR/audio

pushd $DIR

nohup ./stream.webcam.video.sh &
nohup ./stream.webcam.audio.sh &
nohup ./cleanup.dirs.sh &
nohup ./nodejs.server.sh &

popd

